# src/fty-srr.service.in
[Unit]
# In case of recovering a deployment from backups, this does not
# restore EULA acceptance nor license. So still needs (for now)
# to go through the standard wizard + EULA as a prerequisite.
Description=fty-srr service: Save, Restore and Reset
After=malamute.service fty-license-accepted.target
Requires=malamute.service
# Enable FLA... once, do not keep retrying if EULA is not accepted yet
Wants=fty-license-accepted.target
Requisite=fty-license-accepted.target
PartOf=bios.target

[Service]
Type=simple
User=bios
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-srr --config @sysconfdir@/@PACKAGE@/fty-srr.cfg
Restart=always

[Install]
WantedBy=bios.target
