# src/fty-srr.service.in
[Unit]
Description=fty-srr service
After=malamute.service network.target
Requires=malamute.service network.target
PartOf=bios.target

[Service]
Type=simple
User=bios
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-srr --config @sysconfdir@/@PACKAGE@/fty-srr.cfg
Restart=always

[Install]
WantedBy=bios.target
