# src/fty-srr.service.in
[Unit]
# In case of recovering a deployment from backups, this does not
# restore EULA acceptance nor license. So still needs (for now)
# to go through the standard wizard + EULA as a prerequisite.
Description=fty-srr service: Save, Restore and Reset
After=malamute.service network.target fty-license-accepted.target
Requires=malamute.service network.target
# Enable FLA... once, do not keep retrying if EULA is not accepted yet
Wants=fty-license-accepted.target
Requisite=fty-license-accepted.target
PartOf=bios.target

[Service]
Type=simple
User=fty-srr
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-srr --config @sysconfdir@/@PACKAGE@/fty-srr.cfg

[Install]
WantedBy=bios.target
