# src/fty_srr.service.in
# This is a skeleton created by zproject.
# You can add hand-written code here.

[Unit]
Description=fty_srr service
After=malamute.service network.target
Requires=malamute.service network.target
Before=fty_config.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty_srr --config @sysconfdir@/@PACKAGE@/fty_srr.cfg
Restart=always

[Install]
WantedBy=bios.target
